# backend/crontab
# 股票指數自動更新 Cron 配置
# 台北時間 (UTC+8)
# 05:30 — 美股收盤後更新（美東 16:30 = 台北 05:30）
# 09:00 — 早上開盤前確認
# 23:00 — 晚間更新

# 分 時 日 月 週 命令
30 5 * * * cd /app && /usr/local/bin/python /app/update_indices.py >> /var/log/update_indices.log 2>&1
0 9 * * * cd /app && /usr/local/bin/python /app/update_indices.py >> /var/log/update_indices.log 2>&1
0 23 * * * cd /app && /usr/local/bin/python /app/update_indices.py >> /var/log/update_indices.log 2>&1
